# routine/com.wewrite.auto-fetch-kol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wewrite.auto-fetch-kol</string>

  <!-- 每天凌晨 03:00 · 抓 KOL 公众号 RSS · 入 idea_bank -->
  <!-- 错峰原因:用户睡觉时段 · 不抢白天/晚上的 claude token -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/scripts/auto/auto-fetch-kol.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/mahaochen/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/routine/logs/auto-fetch-kol.out.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/routine/logs/auto-fetch-kol.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
