# routine/com.wewrite.auto-discover-kol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wewrite.auto-discover-kol</string>

  <!-- 每天 03:30 · 紧跟 03:05 sync_xhs 之后 · 提供 5 个 AI KOL 候选 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>30</integer>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/scripts/auto/auto-discover-kol.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/mahaochen/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/routine/logs/auto-discover-kol.out.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/routine/logs/auto-discover-kol.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
