# routine/com.wewrite.auto-publish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wewrite.auto-publish</string>

  <!-- 每天 12:00 · auto_publish · 推草稿 + 通知 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>12</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/scripts/auto/auto-publish.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/routine/logs/auto-publish.out.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/routine/logs/auto-publish.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
